# crontab
*/30 * * * * $HOME/.local/bin/sync-flickr >/dev/null 2>&1
